# pi_status.service
[Unit]
Description=Meidesheet bed status
After=network-online.target

[Service]
WorkingDirectory=/home/pi/Scripts/gnb/iot_pi
ExecStart=/bin/sh /home/pi/Scripts/gnb/iot_pi/devStatus.sh
StandardOutput=inherit
StandardError=inherit
Restart=always
User=pi

[Install]
WantedBy=multi-user.target